# deploy/swbfii.service
# Gunicorn process for SWBFII Rankings on a Hetzner VPS.
# Adjust WorkingDirectory/User/paths to match your server layout.

[Unit]
Description=SWBFII Rankings (gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/swbfii
EnvironmentFile=/var/www/swbfii/.env
ExecStart=/var/www/swbfii/venv/bin/gunicorn config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 60 \
    --access-logfile - \
    --error-logfile -
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
